# crontab
# Run every 10 seconds
*/10 * * * * * * python main.py