# crontab
# Run every minute
*/1 * * * * ./main.py